# clkmgr/proxy/clkmgr-proxy.socket.in
# SPDX-License-Identifier: BSD-3-Clause
# SPDX-FileCopyrightText: Copyright © 2025 Intel Corporation.
# See:
# https://www.freedesktop.org/software/systemd/man/systemd.socket.html

[Unit]
Description=Clock Manager service socket
ConditionPathExists=!@sysconfdir@/clkmgr/proxy_not_to_be_run

[Socket]
ListenMessageQueue=/clkmgr
Accept=no
ReceiveBuffer=8K
SendBuffer=8K

[Install]
WantedBy=multi-user.target
Alias=clkmgr-proxy.socket

# vi: set ft=systemd:
